# etcd-server/config/ml-defender-ml-detector.service
# ML Defender — ml-detector
# Install: sudo install -m 0644 ml-defender-ml-detector.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable ml-defender-ml-detector

[Unit]
Description=ML Defender — ml-detector (RandomForest + ONNX inference)
Documentation=https://github.com/alonsoir/argus
After=network.target ml-defender-etcd-server.service ml-defender-rag-ingester.service
Requires=ml-defender-etcd-server.service

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/vagrant/ml-detector/build-active

EnvironmentFile=/etc/ml-defender/build.env
Environment="LD_PRELOAD="
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib"

ExecStart=/vagrant/ml-detector/build-active/ml-detector

Restart=always
RestartSec=5s
StartLimitInterval=300
StartLimitBurst=10

TimeoutStartSec=30s
TimeoutStopSec=15s

# Límites para inferencia ML
LimitNOFILE=32768
MemoryMax=2G
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ml-defender-ml-detector

KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
