# etcd-server/config/ml-defender-etcd-server.service
# ML Defender — etcd-server
# Distribución de configuración, seed distribution, heartbeat registry
# Install: sudo install -m 0644 ml-defender-etcd-server.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable ml-defender-etcd-server

[Unit]
Description=ML Defender — etcd-server (config distribution + seed)
Documentation=https://github.com/alonsoir/argus
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/vagrant

# Seguridad PHASE 3: LD_PRELOAD explícitamente vacío
EnvironmentFile=/etc/ml-defender/build.env
Environment="LD_PRELOAD="
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="LD_LIBRARY_PATH=/usr/local/lib"

ExecStart=/vagrant/etcd-server/build-active/etcd-server

# PHASE 3: Restart=always obligatorio
Restart=always
RestartSec=5s
StartLimitInterval=300
StartLimitBurst=10

TimeoutStartSec=30s
TimeoutStopSec=15s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ml-defender-etcd-server

KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
